# ansible/roles/envoy/templates/envoy.service.j2
# {{ ansible_managed }}
[Unit]
Description=Envoy Proxy Container (Host Network)
Requires=docker.service network-online.target consul.service dnsmasq.service
After=docker.service network-online.target consul.service dnsmasq.service

[Service]
Restart=always
RestartSec=10s

ExecStartPre=-/bin/mkdir -p /etc/envoy
ExecStartPre=-/usr/bin/docker stop envoy-svc
ExecStartPre=-/usr/bin/docker rm envoy-svc
ExecStartPre=-/usr/bin/docker pull {{ envoy_image_repository }}:{{ envoy_image_tag }}

ExecStart=/usr/bin/docker run --rm --network host --name envoy-svc \
  -u root \
  --entrypoint "/usr/local/bin/envoy" \
  -v /etc/envoy/envoy.yaml:/etc/envoy/envoy.yaml:ro \
  -v /etc/letsencrypt:/etc/letsencrypt:ro \
  {{ envoy_image_repository }}:{{ envoy_image_tag }} \
  -c /etc/envoy/envoy.yaml \
  --service-node edge-envoy-instance-1 \
  --service-cluster grewalcc_services
  

ExecStop=/usr/bin/docker stop envoy-svc

[Install]
WantedBy=multi-user.target
