# ansible/roles/envoy/templates/envoy.service.j2
# {{ ansible_managed }} - Template: ~/src/gcc-ansible-wsl/ansible/roles/envoy/templates/envoy.service.j2
[Unit]
Description=Envoy Proxy Container (Host Network)
Requires=docker.service network-online.target consul.service
After=docker.service network-online.target consul.service

[Service]
Restart=always
RestartSec=10s

# Ensure Envoy config directory exists on host (Ansible manages this primarily)
ExecStartPre=-/bin/mkdir -p /etc/envoy

# Stop and remove any existing container with the same name first
ExecStartPre=-/usr/bin/docker stop envoy-svc
ExecStartPre=-/usr/bin/docker rm envoy-svc

# Pull the specified image (optional)
ExecStartPre=-/usr/bin/docker pull {{ envoy_image_repository }}:{{ envoy_image_tag }}

# Run the Envoy container using host networking - All on one line
ExecStart=/usr/bin/docker run --rm --network host --name envoy-svc -v /etc/envoy/envoy.yaml:/etc/envoy/envoy.yaml:ro {{ envoy_image_repository }}:{{ envoy_image_tag }} -c /etc/envoy/envoy.yaml

# Stop the container gracefully
ExecStop=/usr/bin/docker stop envoy-svc

[Install]
WantedBy=multi-user.target
