# ansible/roles/consul/templates/consul.service.j2
[Unit]
Description=Consul Service Discovery Agent (Docker Container)
Requires=docker.service network-online.target
After=docker.service network-online.target

[Service]
Restart=always
RestartSec=10s

# Stop and remove any existing container with the same name first
ExecStartPre=-/usr/bin/docker stop consul-server
ExecStartPre=-/usr/bin/docker rm consul-server

# Run the Consul server container
# Mount the config directory /etc/consul.d read-only.
# Mount the host data directory /opt/consul/data read-write.
# Use the config DIRECTORY argument.
ExecStart=/usr/bin/docker run --rm --network host --name consul-server \
  -v /etc/consul.d:/etc/consul.d:ro \
  -v /opt/consul/data:/consul/data \
  hashicorp/consul:1.20.5 agent -config-dir=/etc/consul.d

# Stop the container gracefully
ExecStop=/usr/bin/docker stop consul-server

[Install]
WantedBy=multi-user.target
